# lnls-epics-archiver.service
[Unit]
Description=LNLS EPICS Archiver Appliances Docker Service
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker-compose -f /usr/local/etc/lnls-epics-archiver/docker-compose.yml up -d
ExecStop=/usr/bin/docker-compose -f /usr/local/etc/lnls-epics-archiver/docker-compose.yml down
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

